# WIFI/配置网络/rc.local
#!/bin/sh -e
# ¼ì²âwifi×´Ì¬
/root/wificonfig.sh &


exit 0
